# launchd/com.usage-tracker.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.usage-tracker.menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/amgad/dev_projects/usage-tracker-minimal/UsageMenuBar/.build/arm64-apple-macosx/release/UsageMenuBar</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/usage-menubar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/usage-menubar.err</string>
</dict>
</plist>
